# deploy/edge/systemd/generated/pix2pi-cert-renew.timer
[Unit]
Description=Run Pix2pi Certbot Renew Foundation

[Timer]
OnCalendar=*-*-* 3:15:00
Persistent=true

[Install]
WantedBy=timers.target
